# ccaa_dir/services/ccaa.service
[Unit]
Description=ccaa
After=network.target

[Service]
Type=forking
ExecStart=/etc/ccaa/ccaa_start.sh

[Install]
WantedBy=multi-user.target