# ccaa_dir/services/ccaa.service
[Unit]
Description=ccaa
After=network.target

[Service]
ExecStart=/usr/sbin/ccaa start

[Install]
WantedBy=multi-user.target